# service/slamd-server/run
#!/bin/sh
exec /usr/local/slamd/start-slamd-server.sh 2>&1
