# service/slamd-server/run
#!/bin/sh
exec /usr/src/slamd/build/package/slamd/start-slamd-server.sh 2>&1
