# service/tmux/run
#!/bin/bash
exec tmux -D -S /tmp/mux/slamd-void 2>&1
